# sparse/etc/systemd/system/systemd-tmpfiles-setup.service.d/execstart.conf
[Service]
ExecStart=
ExecStart=/usr/bin/systemd-tmpfiles --create --remove --boot --exclude-prefix=/dev --exclude-prefix=/sys/fs/cgroup/memory
