# sparse/etc/systemd/system/sensorfwd.service.d/sensors.conf
[Unit]
After=sys-devices-platform-SENSOR.3.device
Requires=sys-devices-platform-SENSOR.3.device
JobTimeoutSec=60 

[Service]
TimeoutStartSec=60
