# sparse/usr/lib/systemd/system/proximity-enable.timer
[Unit]
Description=Enable proximity after booting

[Timer]
OnBootSec=2min

[Install]
WantedBy=multi-user.target


